# dockerfiles/php/runtime/7.2/supervisor/conf.d/crond.ini
[program:crond]
command=crond -f
stdout_logfile=/var/log/supervisor/crond_stdout.log
stdout_logfile_maxbytes=64MB
stderr_logfile=/var/log/supervisor/crond_stderr.log
stderr_logfile_maxbytes=64MB
autorestart=true
startretries=3
