# dockerfiles/php/runtime/7.2/supervisor/conf.d/fpm.ini
[program:fpm]
command=php-fpm7 -F
stdout_logfile=/var/log/supervisor/fpm_stdout.log
stdout_logfile_maxbytes=64MB
stderr_logfile=/var/log/supervisor/fpm_stderr.log
stderr_logfile_maxbytes=64MB
autorestart=true
startretries=3
