# dockerfiles/php/runtime/8.0/supervisor/conf.d/fpm.ini
[program:fpm]
command=php-fpm -F
stdout_logfile=/var/log/supervisord/fpm_stdout.log
stdout_logfile_maxbytes=64MB
stderr_logfile=/var/log/supervisord/fpm_stderr.log
stderr_logfile_maxbytes=64MB
autorestart=true
startretries=3
