# dockerfiles/php/runtime/8.0/supervisor/conf.d/nginx.ini
[program:nginx]
command=nginx -g 'daemon off;'
stdout_logfile=/var/log/supervisord/nginx_stdout.log
stdout_logfile_maxbytes=64MB
stderr_logfile=/var/log/supervisord/nginx_stderr.log
stderr_logfile_maxbytes=64MB
autorestart=true
startretries=3
